# mail/.config/systemd/user/mbsync.service
[Unit]
Description=Run mbsync
After=network.target

[Service]
Type=oneshot
ExecStart=mbsync --all

[Install]
WantedBy=default.target
